# .circleci/resources/etc/supervisor/conf.d/google-chrome.ini
[program:google-chrome.headless]
autostart=true
command=google-chrome
  --headless
  --disable-gpu
  --remote-debugging-port=9222
  --user-data-dir=%(ENV_HOME)s/slash/var/run/google-chrome.headless
  http://localhost/core/misc/druplicon.png
stdout_logfile=%(ENV_HOME)s/slash/var/log/google-chrome.headless.stdout.log
stderr_logfile=%(ENV_HOME)s/slash/var/log/google-chrome.headless.stderr.log
