# .circleci/resources/etc/supervisor/conf.d/php-fpm.ini
[program:php-fpm]
autostart=true
startsecs=10
command=/usr/local/sbin/php-fpm
  -F
  --fpm-config %(ENV_HOME)s/project/.circleci/resources/etc/php-fpm/php-fpm.ini
  -c %(ENV_HOME)s/project/.circleci/resources/etc/php-fpm/php.ini
